# com.footage.claude-skills-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.footage.claude-skills-sync</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/git</string>
    <string>-C</string>
    <string>__HOME__/.claude/skills</string>
    <string>pull</string>
    <string>--ff-only</string>
  </array>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/.claude/logs/skills-sync.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.claude/logs/skills-sync-error.log</string>
</dict>
</plist>
